# .config/systemd/user/acpi-listen.service
[Unit]
Description=ACPI listener

[Service]
Type=simple
ExecStart=/home/juho/bin/acpi-listen

[Install]
WantedBy=default.target
