# .config/systemd/user/docked.target
[Unit]
Description=Laptop Docked
Conflicts=undocked.target
