# .config/systemd/user/dropbox.service
[Unit]
Description=Dropbox

[Install]
WantedBy=default.target

[Service]
Restart=on-failure
Type=simple
ExecStart=/home/juho/bin/dropbox.sh
