# .config/systemd/user/idle-wake.service
[Unit]
Description=Wake Sway

[Service]
Type=simple
ExecStart=/usr/local/bin/swaymsg seat "*" idle_wake switch
