# .config/systemd/user/lxpolkit.service
[Unit]
Description=LXPolkit

[Service]
Type=simple
ExecStart=/usr/bin/lxpolkit
Restart=on-failure
RestartSec=10
KillMode=process
