# .config/systemd/user/pause-spotify.service
[Unit]
Description=Pause Spotify

[Service]
Type=simple
ExecStart=/home/juho/bin/music.sh Pause
