# .config/systemd/user/pause-spotify.service
[Unit]
Description=Pause Spotify

[Service]
Type=simple
ExecStart=/usr/bin/playerctl pause
