# .config/systemd/user/undocked.target
[Unit]
Description=Laptop Undocked
Conflicts=docked.target
