# bin/rc.local
#!/bin/sh
cd /home/juho/bin
sudo -u juho ./pwrmgr.py > /dev/null &
