# bin/rc.local
#!/bin/sh
cd /home/juho/bin
# Run as root because normal user isn't authorized to access Upower without Xsession
./pwrmgr.py | logger -t "pwrmgr" &
